# dot_config/systemd/user/fusuma.service
[Unit]
Description=Starts fusuma user service
After=graphical.target ydotoold.service
Requires=ydotoold.service
ConditionVirtualization=false

[Service]
Type=simple
ExecStart=/usr/local/bin/fusuma
Restart=on-failure

[Install]
WantedBy=multi-user.target

