# dot_config/systemd/user/three-finger-drag.service
[Unit]
Description=Three-finger drag gestures for Linux
After=graphical-session.target

[Service]
Type=exec
ExecStart=%h/.cargo/bin/linux-3-finger-drag
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
