# dot_config/systemd/user/app-copyq@autostart.service.d/override.conf
[Service]
# see: https://github.com/hluk/CopyQ/issues/1243
Environment="QT_QPA_PLATFORM=xcb"
